# forum-crk-maxbot.service
[Unit]
Description=Forum CRK MAX Bot - Бот для регистрации на форуме
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/maxbot/forum_crk_maxbot
Environment="PATH=/root/maxbot/forum_crk_maxbot/venv/bin"
ExecStart=/root/maxbot/forum_crk_maxbot/venv/bin/python /root/maxbot/forum_crk_maxbot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=forum-crk-maxbot

# Защита от перезагрузки при слишком частых перезапусках
StartLimitInterval=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target

